# 08-ansible-03-yandex/playbook/roles/vector/templates/vector.service.j2
[Unit]
Description=Vector Serivce
Documentation=https://vector.dev
After=network-online.target
Requires=network-online.target

[Service]
ExecStart=/usr/bin/vector
ExecReload=/usr/bin/vector validate