# 08-ansible-02-playbook/playbook/roles/vector/templates/vector.service.j2
[Unit]
Description=Vector Serivce

[Service]
ExecStart=/usr/bin/vector
ExecReload=/usr/bin/vector validate